# simple-app.service
[Unit]
Description=Simple App

[Service]
Type=exec
Environment=PATH=%h/child-stdout-demo/build
ExecStart=%h/child-stdout-demo/simple_app/build
